# service/main/finish
#!/bin/sh
cd ../..
exec </dev/null 2>&1
echo finished
